# system/etc/systemd/system/skillfish-gpu-util.service
[Unit]
Description=SkillFishOS GPU utilization sampler (feeds the desktop HUD)
After=multi-user.target

[Service]
ExecStart=/usr/local/bin/skillfish-gpu-util.sh
Restart=always
RestartSec=2
Nice=10

[Install]
WantedBy=multi-user.target
